# service/status-client.service
[Unit]
Description=NodeStatus-Client
Documentation=https://github.com/cokemine/nodestatus-client-go
After=network.target

[Service]
EnvironmentFile=/etc/NodeStatus/config.conf
ExecStart=/usr/local/bin/status-client --dsn $DSN
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target

